# build/php/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/tmp/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:php-fpm]
command=php-fpm
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/php-fpm.log
priority=1

[program:messenger-consume]
command=php /var/www/vhosts/symfony/bin/console messenger:consume async --time-limit=3600
numprocs=1
autorestart=true
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/log/supervisor/messenger-consume.log
priority=100
